# rgbcmd.upstart
# upstart config for rgb
description "rgb"
author "Jason Luther"

start on runlevel [35]
stop on runlevel [!35]

respawn
respawn limit unlimited

script
  /home/pi/git/pi-rgb/rgb
end script